# dot_config/autostart/syncthing.desktop
[Desktop Entry]
Name=SyncThing
Icon=syncthing
Type=Application
Exec=/home/linuxbrew/.linuxbrew/bin/syncthing --no-browser
Terminal=false
